# recommender-api.socket
[Unit]
Description=Socket for Jellyfin Recommender API

[Socket]
ListenStream=8097
Accept=no

[Install]
WantedBy=sockets.target
